# com.ibm.design.checkForCarbonUpdates.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ibm.design.checkForCarbonUpdates</string>
    <key>ProgramArguments</key>
    <array>
        <string>paste-your-path-here/check-for-carbon-update.sh</string>
    </array>
    <key>KeepAlive</key>
    <true />
    <key>ThrottleInterval</key>
    <integer>3600</integer>
</dict>
</plist>